# board/friendlyarm/nanopi-r76s/rootfs_overlay/etc/udev/rules.d/100-wlan-device.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="wlan0", RUN+="/etc/init.d/S100hostapd start"
ACTION=="remove", SUBSYSTEM=="net", KERNEL=="wlan0", RUN+="/etc/init.d/S100hostapd stop"